# .config/systemd/user/parss-desktop-sync.service
[Unit]
Description=PARSS desktop sync (archrice + progs.csv)

[Service]
Type=oneshot
ExecStart=%h/PARSS/scripts/desktop-setup.sh
