# .config/systemd/user/parss-desktop-sync.timer
[Unit]
Description=Run PARSS desktop sync daily

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
